# system/service/lelamp.service
[Unit]
Description=Lelamp Runtime Service
After=network.target

[Service]
Type=simple
User=lelamp
WorkingDirectory=/home/lelamp/lelampv2
ExecStart=/usr/local/bin/uv run main.py console
Restart=always
RestartSec=30
Environment=HOME=/home/lelamp

[Install]
WantedBy=multi-user.target
